# etc/profile.d/10_dist-base-files.sh
#!/bin/sh

## Copyright (C) 2026 - 2026 ENCRYPTED SUPPORT LLC <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

## Indicate to Kicksecure sessions that the profile scripts have been run.
export PROFILE_SCRIPTS_WERE_SOURCED='yes'
